# deployments/etc/systemd/system/markgo.service
[Unit]
Description=MarkGo Engine v1.0.0 - High-Performance Blog Engine
Documentation=https://github.com/vnykmshr/markgo
After=network.target
Wants=network.target

[Service]
Type=simple
User=markgo
Group=markgo
WorkingDirectory=/opt/markgo
ExecStart=/usr/local/bin/markgo
Restart=always
RestartSec=5
TimeoutStopSec=30

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/markgo/logs
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Environment variables  
Environment=ENVIRONMENT=production
Environment=PORT=3000
Environment=BASE_URL=https://yourdomain.com
EnvironmentFile=-/opt/markgo/.env

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=markgo

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
